# tractor_rpi/setup/led-controller.service
[Unit]
Description=LED Status Controller
Documentation=https://github.com/jones2126/tractor2025
After=network.target teensy-bridge.service rtcm-server.service
Wants=teensy-bridge.service rtcm-server.service

[Service]
Type=simple
User=al
WorkingDirectory=/home/al/tractor2025/tractor_rpi
ExecStart=/usr/bin/python3 /home/al/tractor2025/tractor_rpi/led_status_controller.py
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=led-controller

# Resource limits (optional)
# MemoryLimit=256M
# CPUQuota=50%

[Install]
WantedBy=multi-user.target
